# etc/general.crontab
MAILTO=dex-census@lists.alioth.debian.org
@daily cd ~pabs/census ; ./bin/general-cron
