# etc/general.crontab
MAILTO=pabs@debian.org,jonathan@ubuntu.com
@daily cd ~pabs/census ; ./bin/general-cron
