# etc/snapshot.crontab
MAILTO=pabs@debian.org
@daily cd /srv/qa.debian.org/export/derivs/census ; ./bin/snapshot-cron
