# etc/snapshot.crontab
MAILTO=pabs@debian.org,zumbi@kos.to,jonathan@ubuntu.com
@daily cd /srv/qa.debian.org/export/derivs/census ; ./bin/snapshot-cron
