# etc/snapshot.crontab
MAILTO=dex-census@lists.alioth.debian.org
@daily cd /srv/qa.debian.org/export/derivs/census ; ./bin/snapshot-cron
